# files/appliance-initialize.service
[Unit]
Description=Initialize Appliance Database
ConditionPathExists=!/var/lib/pgsql/9.6/data
After=evminit.service memcached.service
Before=evmserverd.service
Requires=memcached.service
[Service]
TimeoutStartSec=10m
Type=oneshot
ExecStart=/bin/appliance-initialize.sh
ExecStartPost=/usr/bin/systemctl disable appliance-initialize
[Install]
WantedBy=multi-user.target
